# deploy/launchd/com.bifrost.client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bifrost.client</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bifrost-client</string>
        <string>-c</string>
        <string>/etc/bifrost/client-config.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/bifrost/client.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/bifrost/client.error.log</string>

    <key>WorkingDirectory</key>
    <string>/var/lib/bifrost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
