# openwrt/etc/init.d/bifrost
#!/bin/sh /etc/rc.common
# Bifrost Proxy Server for OpenWrt
# procd init script
#
# Installation:
#   1. Copy bifrost-server binary to /usr/bin/
#   2. Copy this script to /etc/init.d/bifrost
#   3. Copy config to /etc/config/bifrost (UCI) and/or /etc/bifrost/config.yaml
#   4. Enable: /etc/init.d/bifrost enable
#   5. Start: /etc/init.d/bifrost start

START=90
STOP=10
USE_PROCD=1

PROG=/usr/bin/bifrost-server
PIDFILE=/var/run/bifrost.pid

# Defaults (overridden by UCI config in /etc/config/bifrost).
CONFIG=/etc/bifrost/config.yaml
ENABLED=1
MANAGED=0

# load_uci reads the 'main' section options that control the service itself.
load_uci() {
    . /lib/functions.sh 2>/dev/null || return 0
    config_load bifrost 2>/dev/null || return 0
    config_get ENABLED main enabled "$ENABLED"
    config_get MANAGED main managed "$MANAGED"
    config_get CONFIG main config_file "$CONFIG"
}

# is_true normalises UCI boolean-ish values.
is_true() {
    case "$1" in
        1|yes|true|on|enabled) return 0 ;;
        *) return 1 ;;
    esac
}

# yaml_escape wraps a scalar in double quotes and escapes backslashes and
# double quotes so arbitrary listen addresses / tokens are emitted safely.
yaml_escape() {
    printf '"%s"' "$(printf '%s' "$1" | sed -e 's/\\/\\\\/g' -e 's/"/\\"/g')"
}

# generate_config renders /etc/bifrost/config.yaml from the UCI sections.
# Only the core, statically-known options are mapped; advanced configuration
# (custom backends, routes, auth, cache, TLS, ...) requires managed='0'.
generate_config() {
    local http_listen http_read http_write http_idle http_maxconn
    local socks_listen socks_maxconn
    local log_level log_format log_output
    local alog_enabled alog_format alog_output
    local m_enabled m_listen m_path m_interval
    local web_enabled web_listen
    local api_enabled api_listen api_token api_ws

    config_get http_listen   http   listen        ":8080"
    config_get http_read     http   read_timeout  "30s"
    config_get http_write    http   write_timeout "30s"
    config_get http_idle     http   idle_timeout  "60s"
    config_get http_maxconn  http   max_connections "100"

    config_get socks_listen  socks5 listen        ":1080"
    config_get socks_maxconn socks5 max_connections "100"

    config_get log_level     log    level         "warn"
    config_get log_format    log    format        "text"
    config_get log_output    log    output        "stdout"

    config_get alog_enabled  access_log enabled    "1"
    config_get alog_format   access_log format     "json"
    config_get alog_output   access_log output     "stdout"

    config_get m_enabled     metrics enabled       "0"
    config_get m_listen      metrics listen        ":8090"
    config_get m_path        metrics path          "/metrics"
    config_get m_interval    metrics collection_interval "60s"

    config_get web_enabled   web_ui  enabled       "0"
    config_get web_listen    web_ui  listen        ":8081"

    config_get api_enabled   api     enabled       "1"
    config_get api_listen    api     listen        ":8082"
    config_get api_token     api     token         ""
    config_get api_ws        api     websocket_max_clients "5"

    local dir
    dir=$(dirname "$CONFIG")
    mkdir -p "$dir" || {
        echo "Error: cannot create config directory: $dir" >&2
        return 1
    }

    # Write atomically via a temp file then move into place.
    local tmp
    tmp=$(mktemp "$dir/.bifrost.yaml.XXXXXX" 2>/dev/null) || tmp="$CONFIG.tmp"

    {
        echo "# Generated by /etc/init.d/bifrost from UCI (/etc/config/bifrost)."
        echo "# DO NOT EDIT: regenerated on every start/reload while main.managed='1'."
        echo ""
        echo "server:"
        if [ -n "$http_listen" ]; then
            echo "  http:"
            printf '    listen: %s\n'        "$(yaml_escape "$http_listen")"
            printf '    read_timeout: %s\n'  "$(yaml_escape "$http_read")"
            printf '    write_timeout: %s\n' "$(yaml_escape "$http_write")"
            printf '    idle_timeout: %s\n'  "$(yaml_escape "$http_idle")"
            printf '    max_connections: %s\n' "$http_maxconn"
        fi
        if [ -n "$socks_listen" ]; then
            echo "  socks5:"
            printf '    listen: %s\n'          "$(yaml_escape "$socks_listen")"
            printf '    max_connections: %s\n' "$socks_maxconn"
        fi
        echo "  graceful_period: \"30s\""
        echo ""
        echo "backends:"
        echo "  - name: direct"
        echo "    type: direct"
        echo "    enabled: true"
        echo "    priority: 10"
        echo ""
        echo "routes:"
        echo "  - domains: [\"*\"]"
        echo "    backend: direct"
        echo "    priority: 1"
        echo ""
        echo "auth:"
        echo "  providers:"
        echo "    - name: default"
        echo "      type: none"
        echo "      enabled: true"
        echo "      priority: 1"
        echo ""
        echo "access_log:"
        if is_true "$alog_enabled"; then
            echo "  enabled: true"
        else
            echo "  enabled: false"
        fi
        printf '  format: %s\n' "$(yaml_escape "$alog_format")"
        printf '  output: %s\n' "$(yaml_escape "$alog_output")"
        echo ""
        echo "metrics:"
        if is_true "$m_enabled"; then
            echo "  enabled: true"
        else
            echo "  enabled: false"
        fi
        printf '  listen: %s\n' "$(yaml_escape "$m_listen")"
        printf '  path: %s\n'   "$(yaml_escape "$m_path")"
        printf '  collection_interval: %s\n' "$(yaml_escape "$m_interval")"
        echo ""
        echo "logging:"
        printf '  level: %s\n'  "$(yaml_escape "$log_level")"
        printf '  format: %s\n' "$(yaml_escape "$log_format")"
        printf '  output: %s\n' "$(yaml_escape "$log_output")"
        echo ""
        echo "web_ui:"
        if is_true "$web_enabled"; then
            echo "  enabled: true"
        else
            echo "  enabled: false"
        fi
        printf '  listen: %s\n' "$(yaml_escape "$web_listen")"
        echo ""
        echo "api:"
        if is_true "$api_enabled"; then
            echo "  enabled: true"
        else
            echo "  enabled: false"
        fi
        printf '  listen: %s\n' "$(yaml_escape "$api_listen")"
        if [ -n "$api_token" ]; then
            printf '  token: %s\n' "$(yaml_escape "$api_token")"
        fi
        printf '  websocket_max_clients: %s\n' "$api_ws"
    } > "$tmp" || {
        echo "Error: failed to write generated config to $tmp" >&2
        rm -f "$tmp"
        return 1
    }

    mv -f "$tmp" "$CONFIG" || {
        echo "Error: failed to install generated config at $CONFIG" >&2
        rm -f "$tmp"
        return 1
    }

    # Keep generated config readable only by root: it may contain an API token.
    chmod 600 "$CONFIG" 2>/dev/null

    return 0
}

validate_config() {
    [ -f "$CONFIG" ] || {
        echo "Error: Config file not found: $CONFIG" >&2
        return 1
    }
    [ -x "$PROG" ] || {
        echo "Error: Binary not found or not executable: $PROG" >&2
        return 1
    }
    return 0
}

start_service() {
    load_uci

    case "$ENABLED" in
        0|no|false|off)
            echo "bifrost is disabled in /etc/config/bifrost (option enabled '0')"
            return 0
            ;;
    esac

    if is_true "$MANAGED"; then
        generate_config || {
            echo "Error: refusing to start, config generation failed" >&2
            return 1
        }
    fi

    validate_config || return 1

    procd_open_instance bifrost
    procd_set_param command $PROG start -c $CONFIG

    # Restart on crash: wait 5 seconds, max 5 retries per hour
    procd_set_param respawn 3600 5 5

    # Resource limits for embedded devices
    procd_set_param limits nofile="4096 8192"

    # Log to syslog
    procd_set_param stdout 1
    procd_set_param stderr 1

    # Lower priority to prevent starving other services
    procd_set_param nice 10

    # Store PID
    procd_set_param pidfile $PIDFILE

    procd_close_instance
}

stop_service() {
    # Graceful shutdown via SIGTERM (procd handles this)
    :
}

reload_service() {
    load_uci

    # In managed mode, regenerate the YAML so the reload picks up UCI changes.
    if is_true "$MANAGED"; then
        if ! generate_config; then
            echo "Error: config generation failed, not reloading" >&2
            return 1
        fi
    fi

    local pid
    pid=$(cat "$PIDFILE" 2>/dev/null) || pid=$(pgrep -f "$PROG")
    [ -n "$pid" ] && kill -HUP $pid
}

service_triggers() {
    # Restart/reload when the UCI config (/etc/config/bifrost) changes.
    procd_add_reload_trigger "bifrost"
}

status_service() {
    local pid
    pid=$(cat "$PIDFILE" 2>/dev/null)
    if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
        echo "bifrost is running (PID: $pid)"
        return 0
    else
        echo "bifrost is not running"
        return 1
    fi
}
